# run/systemd/generator/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=local-fs.target

[Mount]
What=/dev/disk/by-partuuid/2ed5c94e-02
Where=/
Type=ext4
Options=defaults,noatime
