# run/systemd/generator/boot.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=local-fs.target
Requires=systemd-fsck@dev-disk-by\x2dpartuuid-2ed5c94e\x2d01.service
After=systemd-fsck@dev-disk-by\x2dpartuuid-2ed5c94e\x2d01.service

[Mount]
What=/dev/disk/by-partuuid/2ed5c94e-01
Where=/boot
Type=vfat
