# var/lib/dpkg/info/dmeventd.list
/.
/lib
/lib/arm-linux-gnueabihf
/lib/arm-linux-gnueabihf/device-mapper
/lib/arm-linux-gnueabihf/device-mapper/libdevmapper-event-lvm2mirror.so
/lib/arm-linux-gnueabihf/device-mapper/libdevmapper-event-lvm2raid.so
/lib/arm-linux-gnueabihf/device-mapper/libdevmapper-event-lvm2snapshot.so
/lib/arm-linux-gnueabihf/device-mapper/libdevmapper-event-lvm2thin.so
/lib/arm-linux-gnueabihf/libdevmapper-event-lvm2.so.2.02
/lib/systemd
/lib/systemd/system
/lib/systemd/system/dm-event.service
/lib/systemd/system/dm-event.socket
/sbin
/sbin/dmeventd
/usr
/usr/share
/usr/share/doc
/usr/share/doc/dmeventd
/usr/share/doc/dmeventd/changelog.Debian.gz
/usr/share/doc/dmeventd/changelog.gz
/usr/share/doc/dmeventd/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/dmeventd.8.gz
/lib/arm-linux-gnueabihf/libdevmapper-event-lvm2mirror.so
/lib/arm-linux-gnueabihf/libdevmapper-event-lvm2raid.so
/lib/arm-linux-gnueabihf/libdevmapper-event-lvm2snapshot.so
/lib/arm-linux-gnueabihf/libdevmapper-event-lvm2thin.so
