# var/lib/dpkg/info/mongodb-server.conffiles
/etc/init.d/mongodb
/etc/init/mongodb.conf
/etc/logrotate.d/mongodb-server
/etc/mongodb.conf
