# var/lib/dpkg/info/mongodb-server.list
/.
/etc
/etc/init
/etc/init/mongodb.conf
/etc/init.d
/etc/init.d/mongodb
/etc/logrotate.d
/etc/logrotate.d/mongodb-server
/etc/mongodb.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/mongodb.service
/usr
/usr/bin
/usr/bin/mongod
/usr/bin/mongos
/usr/share
/usr/share/doc
/usr/share/doc/mongodb-server
/usr/share/doc/mongodb-server/README.Debian
/usr/share/doc/mongodb-server/changelog.Debian.gz
/usr/share/doc/mongodb-server/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/mongodb-server
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/mongod.1.gz
/usr/share/man/man1/mongos.1.gz
/var
/var/lib
/var/lib/mongodb
/var/log
/var/log/mongodb
