# var/lib/dpkg/info/nfs-common.conffiles
/etc/default/nfs-common
/etc/init.d/nfs-common
/etc/request-key.d/id_resolver.conf
