# var/lib/dpkg/info/openvpn.list
/.
/etc
/etc/default
/etc/default/openvpn
/etc/init.d
/etc/init.d/openvpn
/etc/network
/etc/network/if-down.d
/etc/network/if-down.d/openvpn
/etc/network/if-up.d
/etc/network/if-up.d/openvpn
/etc/openvpn
/etc/openvpn/client
/etc/openvpn/server
/etc/openvpn/update-resolv-conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/openvpn-client@.service
/lib/systemd/system/openvpn-server@.service
/lib/systemd/system/openvpn.service
/lib/systemd/system/openvpn@.service
/lib/systemd/system-generators
/lib/systemd/system-generators/openvpn-generator
/usr
/usr/include
/usr/include/openvpn
/usr/include/openvpn/openvpn-msg.h
/usr/include/openvpn/openvpn-plugin.h
/usr/lib
/usr/lib/openvpn
/usr/lib/openvpn/openvpn-plugin-auth-pam.so
/usr/lib/openvpn/openvpn-plugin-down-root.so
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/openvpn.conf
/usr/sbin
/usr/sbin/openvpn
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/openvpn
/usr/share/doc
/usr/share/doc/openvpn
/usr/share/doc/openvpn/AUTHORS
/usr/share/doc/openvpn/COPYING.gz
/usr/share/doc/openvpn/COPYRIGHT.GPL.gz
/usr/share/doc/openvpn/Changes.rst.gz
/usr/share/doc/openvpn/NEWS.Debian.gz
/usr/share/doc/openvpn/PORTS
/usr/share/doc/openvpn/README
/usr/share/doc/openvpn/README.Debian.gz
/usr/share/doc/openvpn/README.IPv6
/usr/share/doc/openvpn/README.auth-pam
/usr/share/doc/openvpn/README.down-root
/usr/share/doc/openvpn/README.polarssl
/usr/share/doc/openvpn/changelog.Debian.gz
/usr/share/doc/openvpn/changelog.gz
/usr/share/doc/openvpn/copyright
/usr/share/doc/openvpn/examples
/usr/share/doc/openvpn/examples/sample-config-files
/usr/share/doc/openvpn/examples/sample-config-files/README
/usr/share/doc/openvpn/examples/sample-config-files/client.conf
/usr/share/doc/openvpn/examples/sample-config-files/firewall.sh
/usr/share/doc/openvpn/examples/sample-config-files/home.up
/usr/share/doc/openvpn/examples/sample-config-files/loopback-client
/usr/share/doc/openvpn/examples/sample-config-files/loopback-server
/usr/share/doc/openvpn/examples/sample-config-files/office.up
/usr/share/doc/openvpn/examples/sample-config-files/openvpn-shutdown.sh
/usr/share/doc/openvpn/examples/sample-config-files/openvpn-startup.sh
/usr/share/doc/openvpn/examples/sample-config-files/server.conf.gz
/usr/share/doc/openvpn/examples/sample-config-files/static-home.conf
/usr/share/doc/openvpn/examples/sample-config-files/static-office.conf
/usr/share/doc/openvpn/examples/sample-config-files/tls-home.conf
/usr/share/doc/openvpn/examples/sample-config-files/tls-office.conf
/usr/share/doc/openvpn/examples/sample-config-files/xinetd-client-config
/usr/share/doc/openvpn/examples/sample-config-files/xinetd-server-config
/usr/share/doc/openvpn/examples/sample-keys
/usr/share/doc/openvpn/examples/sample-keys/README
/usr/share/doc/openvpn/examples/sample-keys/ca.crt
/usr/share/doc/openvpn/examples/sample-keys/ca.key
/usr/share/doc/openvpn/examples/sample-keys/client-ec.crt.gz
/usr/share/doc/openvpn/examples/sample-keys/client-ec.key
/usr/share/doc/openvpn/examples/sample-keys/client-pass.key
/usr/share/doc/openvpn/examples/sample-keys/client.crt.gz
/usr/share/doc/openvpn/examples/sample-keys/client.key
/usr/share/doc/openvpn/examples/sample-keys/client.p12.gz
/usr/share/doc/openvpn/examples/sample-keys/dh2048.pem
/usr/share/doc/openvpn/examples/sample-keys/gen-sample-keys.sh
/usr/share/doc/openvpn/examples/sample-keys/openssl.cnf
/usr/share/doc/openvpn/examples/sample-keys/sample-ca
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/01.pem.gz
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/02.pem.gz
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/03.pem.gz
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/ca.crl
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/ca.crt
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/ca.key
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/client-pass.key
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/client-revoked.crt.gz
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/client-revoked.csr
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/client-revoked.key
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/client.crt.gz
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/client.csr
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/client.key
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/client.p12.gz
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/index.txt
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/index.txt.attr
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/index.txt.attr.old
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/index.txt.old
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/secp256k1.pem
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/serial
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/serial.old
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/server.crt.gz
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/server.csr
/usr/share/doc/openvpn/examples/sample-keys/sample-ca/server.key
/usr/share/doc/openvpn/examples/sample-keys/server-ec.crt.gz
/usr/share/doc/openvpn/examples/sample-keys/server-ec.key
/usr/share/doc/openvpn/examples/sample-keys/server.crt.gz
/usr/share/doc/openvpn/examples/sample-keys/server.key
/usr/share/doc/openvpn/examples/sample-keys/ta.key
/usr/share/doc/openvpn/examples/sample-scripts
/usr/share/doc/openvpn/examples/sample-scripts/auth-pam.pl
/usr/share/doc/openvpn/examples/sample-scripts/bridge-start
/usr/share/doc/openvpn/examples/sample-scripts/bridge-stop
/usr/share/doc/openvpn/examples/sample-scripts/ucn.pl
/usr/share/doc/openvpn/examples/sample-scripts/verify-cn
/usr/share/doc/openvpn/management-notes.txt.gz
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/openvpn.8.gz
/usr/share/openvpn
/usr/share/openvpn/verify-cn
