# var/lib/dpkg/info/pcscd.list
/.
/etc
/etc/init.d
/etc/init.d/pcscd
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/pcscd.service
/lib/systemd/system/pcscd.socket
/usr
/usr/sbin
/usr/sbin/pcscd
/usr/share
/usr/share/doc
/usr/share/doc/pcscd
/usr/share/doc/pcscd/DRIVERS
/usr/share/doc/pcscd/README.gz
/usr/share/doc/pcscd/SECURITY
/usr/share/doc/pcscd/changelog.Debian.gz
/usr/share/doc/pcscd/changelog.gz
/usr/share/doc/pcscd/copyright
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/reader.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/pcscd.8.gz
