# var/lib/dpkg/info/raspberrypi-sys-mods.list
/.
/etc
/etc/modprobe.d
/etc/modprobe.d/blacklist-rtl8192cu.conf
/etc/modprobe.d/blacklist-rtl8xxxu.conf
/etc/sudoers.d
/etc/sudoers.d/010_pi-nopasswd
/etc/sysctl.d
/etc/sysctl.d/98-rpi.conf
/etc/udev
/etc/udev/rules.d
/etc/udev/rules.d/99-com.rules
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/apply_noobs_os_config.service
/lib/systemd/system/regenerate_ssh_host_keys.service
/lib/systemd/system/rpi-display-backlight.service
/lib/systemd/system/sshswitch.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/15-i2c-modprobe.rules
/usr
/usr/lib
/usr/lib/raspberrypi-sys-mods
/usr/lib/raspberrypi-sys-mods/i2cprobe
/usr/share
/usr/share/doc
/usr/share/doc/raspberrypi-sys-mods
/usr/share/doc/raspberrypi-sys-mods/NEWS.Debian.gz
/usr/share/doc/raspberrypi-sys-mods/changelog.gz
/usr/share/doc/raspberrypi-sys-mods/copyright
